# systemd-user/classes.path
[Path]
Unit=classes.service
PathChanged=%h/sync
# trigger on changes to a file not just create/delete

[Install]
WantedBy=default.target