# systemd-user/fusuma.service
[Unit]
Description=start fusuma
Documentation=man:fusuma

[Service]
Type=notify
ExecStart=fusuma

[Install]
WantedBy=default.target
